# services/hf256.service
[Unit]
Description=HF-256 Messaging Application
After=freedvtnc2.service
ConditionPathExists=/etc/hf256/.setup_complete
ConditionPathExists=/etc/hf256/settings.json

[Service]
Type=simple
User=root
WorkingDirectory=/opt/hf256
ExecStart=/usr/bin/python3 /opt/hf256/hf256/main.py
Restart=on-failure
RestartSec=15
StandardOutput=journal
StandardError=journal
Environment=PYTHONUNBUFFERED=1
Environment=PYTHONPATH=/opt/hf256

[Install]
WantedBy=multi-user.target
