# services/rigctld.service
[Unit]
Description=Hamlib rigctld Radio CAT Control
After=hf256-wlan.service
ConditionPathExists=/etc/hf256/config.env

[Service]
Type=simple
User=pi
EnvironmentFile=/etc/hf256/config.env
ExecStart=/bin/bash -c "${RIGCTLD_CMD}"
Restart=on-failure
RestartSec=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
